# mock/supervisord.ini
[supervisord]
nodaemon=true

[program:mysqld]
command=/usr/bin/pidproxy /tmp/mysqld.pid /usr/bin/mysqld_safe --user=root --pid-file=/tmp/mysqld.pid
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:pprotein-mock]
command=/usr/local/bin/pprotein-mock
autostart=true
autorestart=true

[program:pprotein-agent]
command=/usr/local/bin/pprotein-agent
autostart=true
autorestart=true
